# rc.local
#!/bin/bash


if [[ -f /boot/config.txt ]]; then
    # root user?
    apt-get update
    apt-get install python3-pip -y
    pip3 install tplink-cloud-api
    pip3 install RPi.GPIO
    rm /boot/config.txt
done

